# lib/systemd/system/openfortivpn@.service.in
[Unit]
Description=OpenFortiVPN for %I
After=network-online.target
Wants=network-online.target systemd-networkd-wait-online.service
Documentation=man:openfortivpn(1)
Documentation=https://github.com/adrienverge/openfortivpn#readme
Documentation=https://github.com/adrienverge/openfortivpn/wiki

[Service]
Type=notify
PrivateTmp=true
ExecStart=@BINDIR@/openfortivpn -c @SYSCONFDIR@/openfortivpn/%I
Restart=on-failure
OOMScoreAdjust=-100

[Install]
WantedBy=multi-user.target
